# paradata.service
[Unit]
Description=paradata.fr - Carte collaborative des zones sensibles
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/paradata
Environment="PATH=/var/www/paradata/.venv/bin"
ExecStart=/var/www/paradata/.venv/bin/uvicorn app.main:socket_app --host 127.0.0.1 --port 8000 --workers 4
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
